# etc/deploy/events-view.conf
description "yirga node process"

start on runlevel [2]
stop on runlevel [016]

console log
chdir "/var/www/yirga"
env NODE_ENV=production
exec npm start >> /var/log/yirga.log 2>&1
respawn
respawn limit 5 15
